# root/etc/s6-overlay/s6-rc.d/init-mod-universal-wireguard-prometheus-exporter-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-universal-wireguard-prometheus-exporter-add-package/run